# host_setup/pi.crontab
######################################################################################################################
# pi.crontab
##############
# m h  dom mon dow   command
######################################################################################################################
## Update Docker images and re-provision containers
00 04 * * * /mnt/appdata/stop_update_start.sh > /mnt/appdata/stop_update_start.cronlog 2>&1
